# nsebot.service
[Unit]
Description=Market Pulse v1.0 Intelligence Engine
After=network.target

[Service]
Type=simple
User=root
# Change this path to your actual project directory on Ubuntu
WorkingDirectory=/root/nse2
# Ensure .env is loaded
EnvironmentFile=/root/nse2/nse_monitor/.env
# Use system python or venv python
ExecStart=/usr/bin/python3 -m nse_monitor.main

# Self-Healing: Restart on failure AND normal exit (os._exit(0))
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
